# akile-monitor.service
[Unit]
Description=AKILE Stock Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/akile-monitor
ExecStart=/usr/bin/python3 -u monitor.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
